# wega-hpg/scripts/wega-hpg.service
[Unit]
After=network.target

[Service]
ExecStart=/var/WEGA/wega-hpg/entrypoint.sh

[Install]
WantedBy=default.target